# deskai.service
[Unit]
Description=DeskAI Background Service
After=network.target ollama.service

[Service]
Type=simple

# CRITICAL
WorkingDirectory=%h/Downloads/deskai

# Absolute python from venv
ExecStart=%h/Downloads/deskai/desk/bin/python main.py

# Environment is NOT inherited
Environment=PATH=%h/Downloads/deskai/desk/bin:/usr/local/bin:/usr/bin
Environment=HOME=%h

Restart=always
RestartSec=3

# Avoid restart storms
StartLimitIntervalSec=60
StartLimitBurst=5

[Install]
WantedBy=default.target

